# compliance/mappings/cis_controls_v8.yaml
# CIS Controls v8 Mappings
# Format: technique_id → list of control mappings
# Reference: https://www.cisecurity.org/controls/v8

T1059.001:
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "PowerShell Script Block Logging (EID 4104) captures encoded command content for audit"
  - control_id: "8.8"
    control_name: "Collect Command-Line Audit Logs"
    mapping_note: "Process creation logging (EID 4688) captures -EncodedCommand flag in PowerShell invocations"
  - control_id: "10.1"
    control_name: "Deploy and Maintain Anti-Malware Software"
    mapping_note: "AMSI integration with AV/EDR products provides secondary detection layer for encoded scripts"
  - control_id: "4.8"
    control_name: "Uninstall or Disable Unnecessary Services on Enterprise Assets and Software"
    mapping_note: "Constrained Language Mode and execution policy restrict unauthorized PowerShell usage"

T1053.005:
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "Task Scheduler operational log (EID 4698) records scheduled task creation events"
  - control_id: "5.4"
    control_name: "Restrict Administrator Privileges to Dedicated Administrator Accounts"
    mapping_note: "Scheduled task creation by non-admin accounts is anomalous and should be restricted"
  - control_id: "4.1"
    control_name: "Establish and Maintain a Secure Configuration Process"
    mapping_note: "Baseline of legitimate scheduled tasks enables detection of unauthorized additions"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Sysmon EID 1 captures schtasks.exe execution with suspicious arguments"

T1547.001:
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "Sysmon EID 13 (RegSetValue) logs writes to HKCU/HKLM Run keys"
  - control_id: "4.1"
    control_name: "Establish and Maintain a Secure Configuration Process"
    mapping_note: "Baseline of legitimate Run key entries enables detection of unauthorized persistence"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Registry monitoring detects modification of autorun keys used for persistence"
  - control_id: "10.2"
    control_name: "Configure Automatic Anti-Malware Signature Updates"
    mapping_note: "AV products scanning registry persistence locations provide secondary detection"

T1078:
  - control_id: "5.1"
    control_name: "Establish and Maintain an Inventory of Accounts"
    mapping_note: "Account inventory enables detection of unexpected or dormant accounts being used"
  - control_id: "5.3"
    control_name: "Disable Dormant Accounts"
    mapping_note: "Dormant accounts used by attackers are detectable through logon anomaly monitoring"
  - control_id: "5.4"
    control_name: "Restrict Administrator Privileges to Dedicated Administrator Accounts"
    mapping_note: "Admin logon from unexpected hosts (EID 4624 logon type 3) triggers detection"
  - control_id: "6.3"
    control_name: "Require MFA for Externally-Exposed Applications"
    mapping_note: "MFA prevents credential reuse even when valid credentials are compromised"
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "Logon events (EID 4624, 4625, 4648) provide visibility into credential usage patterns"

T1003.001:
  - control_id: "5.4"
    control_name: "Restrict Administrator Privileges to Dedicated Administrator Accounts"
    mapping_note: "LSASS access requires elevated privileges; restricting admin rights limits blast radius"
  - control_id: "10.5"
    control_name: "Enable Anti-Exploitation Features"
    mapping_note: "LSA Protection (RunAsPPL) and Credential Guard prevent unauthorized LSASS memory reads"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Sysmon EID 10 (ProcessAccess) on lsass.exe is the primary detection signal"
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "EID 4656 (object handle request) provides audit trail of LSASS access attempts"
  - control_id: "5.2"
    control_name: "Use Unique Passwords"
    mapping_note: "Unique credentials limit lateral movement impact after LSASS dump"

T1110.001:
  - control_id: "5.2"
    control_name: "Use Unique Passwords"
    mapping_note: "Strong unique passwords increase resistance to password spraying attacks"
  - control_id: "6.3"
    control_name: "Require MFA for Externally-Exposed Applications"
    mapping_note: "MFA renders password spraying ineffective even when passwords are guessed correctly"
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "EID 4625 (failed logon) and EID 4771 (Kerberos failure) provide spray detection signals"
  - control_id: "6.7"
    control_name: "Use DNS Filtering Services"
    mapping_note: "DNS filtering can block C2 infrastructure used after successful spray"
  - control_id: "4.7"
    control_name: "Manage Default Accounts on Enterprise Assets and Software"
    mapping_note: "Default or common account names are primary targets for password spraying"

T1055.001:
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Sysmon EID 8 (CreateRemoteThread) detects cross-process injection attempts"
  - control_id: "10.5"
    control_name: "Enable Anti-Exploitation Features"
    mapping_note: "Windows Defender Exploit Guard process injection mitigations block DLL injection"
  - control_id: "8.8"
    control_name: "Collect Command-Line Audit Logs"
    mapping_note: "Process creation logging captures injection loader execution"
  - control_id: "10.1"
    control_name: "Deploy and Maintain Anti-Malware Software"
    mapping_note: "AV/EDR solutions monitor for remote thread creation and memory write patterns"

T1218.011:
  - control_id: "2.5"
    control_name: "Allowlist Authorized Software"
    mapping_note: "Application allowlisting blocks unauthorized DLLs executed via rundll32.exe"
  - control_id: "8.8"
    control_name: "Collect Command-Line Audit Logs"
    mapping_note: "Rundll32 with unusual arguments captured in EID 4688 / Sysmon EID 1"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Sysmon EID 3 detects network connections initiated by rundll32.exe (high confidence)"
  - control_id: "4.8"
    control_name: "Uninstall or Disable Unnecessary Services"
    mapping_note: "Restrict execution of LOLBINs through AppLocker or WDAC policies"

T1083:
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "File access auditing (EID 4663) with SACLs captures unauthorized directory enumeration"
  - control_id: "3.3"
    control_name: "Configure Data Access Control Lists"
    mapping_note: "Proper ACLs limit what an attacker can enumerate after initial access"
  - control_id: "8.8"
    control_name: "Collect Command-Line Audit Logs"
    mapping_note: "dir, tree, Get-ChildItem execution captured via process creation logging"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Sysmon EID 1 detects high-frequency discovery command execution patterns"

T1082:
  - control_id: "8.8"
    control_name: "Collect Command-Line Audit Logs"
    mapping_note: "systeminfo.exe, whoami.exe, net.exe execution captured in process creation logs"
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "Audit logs capture system discovery activity enabling timeline reconstruction"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Multiple discovery commands in rapid succession from same parent is high-confidence signal"
  - control_id: "4.1"
    control_name: "Establish and Maintain a Secure Configuration Process"
    mapping_note: "Baseline of normal discovery activity helps tune detection thresholds"

T1021.002:
  - control_id: "12.2"
    control_name: "Establish and Maintain a Secure Network Architecture"
    mapping_note: "Network segmentation limits SMB lateral movement to segmented zones"
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "EID 5140 (share access) and EID 5145 (share object check) log SMB admin share usage"
  - control_id: "5.4"
    control_name: "Restrict Administrator Privileges to Dedicated Administrator Accounts"
    mapping_note: "Admin share access (C$, ADMIN$) requires admin rights; restricting these limits lateral movement"
  - control_id: "12.8"
    control_name: "Establish and Maintain Dedicated Computing Resources for All Administrative Work"
    mapping_note: "Privileged Access Workstations limit which hosts can legitimately access admin shares"
  - control_id: "13.4"
    control_name: "Perform Traffic Filtering Between Network Segments"
    mapping_note: "Blocking SMB (port 445) between workstations prevents peer-to-peer lateral movement"

T1070.001:
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "EID 1102 (Security log cleared) and EID 104 (System log cleared) fire at moment of clearing"
  - control_id: "8.3"
    control_name: "Ensure Adequate Audit Log Storage"
    mapping_note: "Insufficient storage causes automatic log overwrite, mimicking attacker log clearing"
  - control_id: "8.9"
    control_name: "Centralize Audit Logs"
    mapping_note: "Forwarding logs to a central SIEM/collector before they can be cleared is the primary mitigation"
  - control_id: "5.4"
    control_name: "Restrict Administrator Privileges to Dedicated Administrator Accounts"
    mapping_note: "Only admins can clear Security logs; restricting admin rights limits who can destroy evidence"

T1486:
  - control_id: "11.1"
    control_name: "Establish and Maintain a Data Recovery Practice"
    mapping_note: "Offline/offsite backups are the primary recovery mechanism after ransomware encryption"
  - control_id: "11.2"
    control_name: "Perform Automated Backups"
    mapping_note: "Automated backups reduce recovery point objective (RPO) after encryption impact"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Sysmon EID 11 (high-volume FileCreate with unusual extensions) and EID 23 (FileDelete) indicate encryption"
  - control_id: "8.2"
    control_name: "Collect Audit Logs"
    mapping_note: "EID 4663 mass file access events indicate active encryption activity"
  - control_id: "3.3"
    control_name: "Configure Data Access Control Lists"
    mapping_note: "Least privilege file permissions limit the blast radius of file encryption"

T1490:
  - control_id: "11.1"
    control_name: "Establish and Maintain a Data Recovery Practice"
    mapping_note: "VSS deletion directly targets Windows backup recovery — offline backups are essential"
  - control_id: "11.3"
    control_name: "Protect Recovery Data"
    mapping_note: "Protecting VSS and backup infrastructure from deletion is a critical hardening step"
  - control_id: "8.8"
    control_name: "Collect Command-Line Audit Logs"
    mapping_note: "vssadmin delete shadows and wmic shadowcopy delete captured via process creation logging"
  - control_id: "5.4"
    control_name: "Restrict Administrator Privileges to Dedicated Administrator Accounts"
    mapping_note: "VSS deletion requires elevated privileges; admin restrictions limit who can inhibit recovery"

T1566.001:
  - control_id: "9.1"
    control_name: "Ensure Use of Only Fully Supported Browsers and Email Clients"
    mapping_note: "Up-to-date Office applications have macro security controls that limit attachment execution"
  - control_id: "9.6"
    control_name: "Block Unnecessary File Types"
    mapping_note: "Email gateway filtering of macro-enabled attachments reduces phishing delivery success"
  - control_id: "14.1"
    control_name: "Establish and Maintain a Security Awareness Program"
    mapping_note: "User awareness training reduces click rates on spearphishing attachments"
  - control_id: "8.8"
    control_name: "Collect Command-Line Audit Logs"
    mapping_note: "Office process spawning cmd.exe or powershell.exe captured via process creation logging"
  - control_id: "13.7"
    control_name: "Deploy a Host-Based Intrusion Detection Solution"
    mapping_note: "Sysmon EID 1 parent-child chain (winword.exe → powershell.exe) is near-zero FP detection"