# compliance/mappings/iso27001.yaml
# ISO 27001:2022 Mappings
# Format: technique_id → list of control mappings
# Reference: ISO/IEC 27001:2022 Annex A Controls
# Note: Control IDs follow ISO 27001:2022 numbering (5.x–8.x)
# Legacy 2013 IDs (A.x.x.x) noted where commonly referenced

T1059.001:
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "PowerShell Script Block Logging (EID 4104) constitutes event logging of user and system activity"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Monitoring of encoded PowerShell execution detects potential malicious command obfuscation"
  - control_id: "8.7"
    control_name: "Protection Against Malware"
    mapping_note: "AMSI and endpoint controls provide malware protection layer against malicious scripts"
  - control_id: "5.25"
    control_name: "Assessment and Decision on Information Security Events"
    mapping_note: "Encoded PowerShell alert requires assessment to determine if event constitutes an incident"

T1053.005:
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "Task Scheduler operational log (EID 4698) records scheduled task creation for audit purposes"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Monitoring for unauthorized scheduled task creation detects persistence establishment"
  - control_id: "5.37"
    control_name: "Documented Operating Procedures"
    mapping_note: "Documented baseline of authorized scheduled tasks enables change detection"
  - control_id: "8.18"
    control_name: "Use of Privileged Utility Programs"
    mapping_note: "schtasks.exe is a privileged utility whose usage must be monitored and controlled"

T1547.001:
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "Sysmon EID 13 (RegSetValue) provides log record of registry modification activity"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Registry autorun key monitoring detects unauthorized persistence mechanisms"
  - control_id: "8.9"
    control_name: "Configuration Management"
    mapping_note: "Baseline configuration of registry Run keys enables detection of unauthorized additions"
  - control_id: "5.25"
    control_name: "Assessment and Decision on Information Security Events"
    mapping_note: "Unexpected Run key modification triggers event assessment workflow"

T1078:
  - control_id: "5.16"
    control_name: "Identity Management"
    mapping_note: "Identity management controls govern legitimate account usage; anomalous usage is detectable"
  - control_id: "5.17"
    control_name: "Authentication Information"
    mapping_note: "Valid account abuse exploits authentication credentials; strong authentication policies mitigate"
  - control_id: "5.18"
    control_name: "Access Rights"
    mapping_note: "Access rights reviews identify over-privileged accounts susceptible to abuse"
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "Logon events (EID 4624, 4625, 4648, 4672) provide audit trail of credential usage"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Logon anomaly monitoring (unusual hours, new hosts) detects valid account misuse"

T1003.001:
  - control_id: "5.17"
    control_name: "Authentication Information"
    mapping_note: "LSASS dump directly compromises authentication information (hashes, tickets, passwords)"
  - control_id: "5.18"
    control_name: "Access Rights"
    mapping_note: "Restricting access rights to LSASS process limits who can perform credential dumping"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Sysmon EID 10 (ProcessAccess on lsass.exe) is the primary monitoring signal"
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "EID 4656 records object handle requests to LSASS for audit purposes"
  - control_id: "5.26"
    control_name: "Response to Information Security Incidents"
    mapping_note: "LSASS access is a high-severity incident requiring immediate credential rotation response"

T1110.001:
  - control_id: "5.17"
    control_name: "Authentication Information"
    mapping_note: "Password spraying targets authentication information; password quality policies mitigate"
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "EID 4625 (failed logon) and EID 4771 (Kerberos pre-auth failure) provide detection data"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Threshold-based monitoring of failed logons across multiple accounts detects spraying"
  - control_id: "5.16"
    control_name: "Identity Management"
    mapping_note: "Account lockout policies and MFA controls reduce effectiveness of password spraying"

T1055.001:
  - control_id: "8.7"
    control_name: "Protection Against Malware"
    mapping_note: "Endpoint protection solutions monitor for DLL injection and remote thread creation patterns"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Sysmon EID 8 (CreateRemoteThread) and EID 10 (ProcessAccess with write) are key signals"
  - control_id: "8.22"
    control_name: "Filtering of Web Services"
    mapping_note: "Web filtering limits delivery of injection payloads from external sources"
  - control_id: "5.25"
    control_name: "Assessment and Decision on Information Security Events"
    mapping_note: "Cross-process injection events require immediate assessment as likely active compromise"

T1218.011:
  - control_id: "8.19"
    control_name: "Installation of Software on Operational Systems"
    mapping_note: "Application control policies restrict unauthorized DLL execution via rundll32.exe"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Monitoring rundll32.exe for network connections (EID 3) is near-zero false positive"
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "Process creation logging captures rundll32 LOLBIN abuse in command-line arguments"
  - control_id: "8.9"
    control_name: "Configuration Management"
    mapping_note: "AppLocker and WDAC policies enforce authorized use of system utilities"

T1083:
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "File access auditing (EID 4663) with configured SACLs logs directory enumeration"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "High-frequency file and directory access events indicate active discovery reconnaissance"
  - control_id: "5.12"
    control_name: "Classification of Information"
    mapping_note: "Information classification enables identification of sensitive paths being enumerated"
  - control_id: "5.15"
    control_name: "Access Control"
    mapping_note: "Least privilege access controls limit the scope of what an attacker can discover"

T1082:
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "Process creation logs capture systeminfo.exe, whoami.exe, and net.exe execution"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Cluster of discovery commands from single parent process is high-confidence compromise indicator"
  - control_id: "5.37"
    control_name: "Documented Operating Procedures"
    mapping_note: "Documented baseline of normal admin tool usage helps tune discovery detection thresholds"
  - control_id: "8.9"
    control_name: "Configuration Management"
    mapping_note: "System configuration documentation enables detection of unauthorized enumeration activity"

T1021.002:
  - control_id: "8.20"
    control_name: "Networks Security"
    mapping_note: "Network segmentation and SMB port filtering (445) limit lateral movement via admin shares"
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "EID 5140 (network share object accessed) and EID 5145 log admin share access events"
  - control_id: "5.15"
    control_name: "Access Control"
    mapping_note: "Restricting admin share access (C$, ADMIN$) to authorized management systems only"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Admin share access from non-management workstations is anomalous and monitored"
  - control_id: "8.21"
    control_name: "Security of Network Services"
    mapping_note: "Network service controls govern which hosts may initiate SMB connections to sensitive systems"

T1070.001:
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "EID 1102 (Security log cleared) and EID 104 (System log cleared) are generated at time of clearing"
  - control_id: "8.17"
    control_name: "Clock Synchronization"
    mapping_note: "Accurate timestamps on log entries are critical for reconstructing timelines after log clearing"
  - control_id: "5.33"
    control_name: "Protection of Records"
    mapping_note: "Log records must be protected from unauthorized modification or deletion per ISO 27001"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Any occurrence of EID 1102/104 triggers immediate monitoring escalation"
  - control_id: "5.28"
    control_name: "Collection of Evidence"
    mapping_note: "Centralized log forwarding ensures evidence is preserved even when local logs are cleared"

T1486:
  - control_id: "8.13"
    control_name: "Information Backup"
    mapping_note: "Regular offline backups are the primary control for recovery from ransomware encryption"
  - control_id: "5.30"
    control_name: "ICT Readiness for Business Continuity"
    mapping_note: "Business continuity planning must account for ransomware as a primary disruption scenario"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "High-volume FileCreate (EID 11) and FileDelete (EID 23) events signal active encryption"
  - control_id: "5.26"
    control_name: "Response to Information Security Incidents"
    mapping_note: "Ransomware encryption requires immediate IR response per documented incident procedures"
  - control_id: "5.29"
    control_name: "Information Security During Disruption"
    mapping_note: "Controls must remain operational during ransomware disruption to support recovery"

T1490:
  - control_id: "8.13"
    control_name: "Information Backup"
    mapping_note: "VSS (Volume Shadow Copy) deletion directly destroys Windows backup recovery capability"
  - control_id: "5.30"
    control_name: "ICT Readiness for Business Continuity"
    mapping_note: "Business continuity plans must include recovery paths that do not rely solely on VSS"
  - control_id: "8.15"
    control_name: "Logging"
    mapping_note: "vssadmin.exe and wmic.exe execution captured in process creation logs (EID 4688 / Sysmon EID 1)"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Any execution of vssadmin delete shadows is a critical alert requiring immediate response"
  - control_id: "5.33"
    control_name: "Protection of Records"
    mapping_note: "Backup records and recovery data must be protected from deletion by unauthorized processes"

T1566.001:
  - control_id: "8.23"
    control_name: "Web Filtering"
    mapping_note: "Web and email filtering controls reduce delivery success rate of malicious attachments"
  - control_id: "6.3"
    control_name: "Information Security Awareness, Education and Training"
    mapping_note: "User awareness training reduces susceptibility to spearphishing attachment social engineering"
  - control_id: "8.7"
    control_name: "Protection Against Malware"
    mapping_note: "Email gateway AV scanning and Office macro controls provide pre-execution protection"
  - control_id: "8.16"
    control_name: "Monitoring Activities"
    mapping_note: "Office application spawning shell processes (EID 1 parent-child chain) is primary detection"
  - control_id: "5.25"
    control_name: "Assessment and Decision on Information Security Events"
    mapping_note: "Phishing alert requires assessment to determine if successful execution occurred"
  - control_id: "8.8"
    control_name: "Management of Technical Vulnerabilities"
    mapping_note: "Keeping Office applications patched reduces exploitation of attachment-borne vulnerabilities"